# srvconf/couchdb-service-manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'> 
<service_bundle type='manifest' name='export'>
  <service name='application/database/couch' type='service' version='0'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    <dependency name='filesystem-local' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
  
  <!--
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name='couchdb-epmd' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/pkgsrc/epmd' />
    </dependency>
  -->
    <method_context>
      <method_credential user='couchdb' group='couchdb'/>
      <method_environment>
      <envvar name="PATH" value="/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin:/bin" />
      <envvar name="HOME" value="/var/db/couchdb" />
      </method_environment>
    </method_context>
    <exec_method name='start' type='method' exec='/opt/local/bin/couchdb -b ' timeout_seconds='300'/>
    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="60" />
    <stability value='Evolving'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Apache CouchDB</loctext>
      </common_name>
      <documentation>
        <manpage title='couchdb' section='1M'/>
        <doc_link name='incubator.apache.org' uri='http://incubator.apache.org/couchdb/'/>
      </documentation>
    </template>
  </service>
 
</service_bundle>